# hdfs-xrootd-healer/config/hdfs-xrootd-healer.cron
SYS_LOCK="/var/lock/subsys/hdfs-xrootd-healer"
CRON_LOCK="/var/lock/hdfs-xrootd-healer/hdfs-xrootd-healer-cron"
HEALER_SCRIPT="/usr/sbin/hdfs-xrootd-healer"

0 4 * * * root [ ! -f $SYS_LOCK ] || [ -f $CRON_LOCK ] || { touch $CRON_LOCK && $HEALER_SCRIPT;rm $CRON_LOCK; }
